# report/tools/service_files/test6/gtrm4.service
[Unit]

[Service]
ExecStart=/mnt/flash/test-gtrm-app 50 20 0 0 0 0.01 0.2 0.04
CPUShares=100
Slice=gtrm.slice

