# report/tools/service_files/test4/gtrm3.service
[Unit]

[Service]
ExecStart=/mnt/flash/test-gtrm-app 75 10 0 0 0 0.01 1 0.04
CPUShares=100
Slice=gtrm.slice

