# report/tools/service_files/test4/gtrm2.service
[Unit]

[Service]
ExecStart=/mnt/flash/test-gtrm-app 150 20 0 0 0 0.6 1 0.04
CPUShares=100
Slice=gtrm.slice

