# report/tools/service_files/test5/gtrm4.service
[Unit]

[Service]
ExecStart=/mnt/flash/test-gtrm-app 75 0 30 0 0 0.01 0.8 0.04
CPUShares=300
Slice=gtrm.slice

